# deploy/approach-a/systemd/jupyter@.service
[Unit]
Description=JupyterLab for %i
After=network.target

[Service]
User=%i
Group=%i
Environment=HOME=/home/%i
EnvironmentFile=/etc/pegasus-studio/users/%i.env
WorkingDirectory=/home/%i
ExecStartPre=/usr/bin/install -d /home/%i/work
# All server settings (port from STUDIO_JUPYTER_PORT, no-auth-behind-gateway,
# allow_remote_access, etc.) live in the shared config file.
ExecStart=/opt/pegasus-studio/jupyter/.venv/bin/jupyter-lab \
    --config=/etc/pegasus-studio/jupyter_server_config.py
Restart=on-failure
RestartSec=3
MemoryMax=2G

[Install]
WantedBy=multi-user.target
