# deploy/approach-a/systemd/studio-reaper.service
[Unit]
Description=PegasusAI Studio session reaper (workflow-aware idle stop)

[Service]
Type=oneshot
ExecStart=/opt/pegasus-studio/bin/studio-reaper.sh
